# com.weather.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.weather.daily</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/miniconda3/bin/python3</string>
        <string>/Users/incarnation/Desktop/准动（银豹）/气象预测货物销售状况/每日自动执行.py</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/incarnation/Desktop/准动（银豹）/气象预测货物销售状况</string>
    
    <key>StandardOutPath</key>
    <string>/Users/incarnation/Desktop/准动（银豹）/气象预测货物销售状况/logs/daily_weather_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/incarnation/Desktop/准动（银豹）/气象预测货物销售状况/logs/daily_weather_stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/incarnation/Desktop/准动（银豹）/气象预测货物销售状况</string>
    </dict>
</dict>
</plist> 